# deploy/flowweave-vpn-client.service
[Unit]
Description=FlowWeave TUN VPN client
Documentation=file:/usr/local/share/doc/flowweave/VPN_RESEARCH.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=notify
NotifyAccess=main
User=flowweave
Group=flowweave
UMask=0077
RuntimeDirectory=flowweave-vpn-client
RuntimeDirectoryMode=0700
Environment=FLOWWEAVE_VPN_CLIENT_RELOAD_SOCKET=/run/flowweave-vpn-client/reload.sock

ExecStartPre=+/usr/local/bin/flowweave-vpn-net prepare-client /etc/flowweave/vpn-client.json /run/flowweave-vpn-client.network.json @flowweave
ExecStart=/usr/local/bin/flowweave-vpn client /etc/flowweave/vpn-client.json
ExecStartPost=+/usr/local/bin/flowweave-vpn-net activate-client /etc/flowweave/vpn-client.json /run/flowweave-vpn-client.network.json
ExecReload=/usr/local/bin/flowweave-vpn reload-client /run/flowweave-vpn-client/reload.sock
ExecStopPost=+/usr/local/bin/flowweave-vpn-net deactivate-client /run/flowweave-vpn-client.network.json
ExecStopPost=+/usr/local/bin/flowweave-vpn-net cleanup /run/flowweave-vpn-client.network.json

Restart=on-failure
RestartSec=2s
TimeoutStartSec=90s
TimeoutStopSec=30s
KillMode=control-group
OOMPolicy=stop
LimitNOFILE=65536
TasksMax=512
MemoryMax=1G
IPAccounting=true

CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=true
DevicePolicy=closed
DeviceAllow=/dev/net/tun rw
PrivateDevices=false
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
