# deploy/public-soak/flowweave-public-soak-client.service
[Unit]
Description=FlowWeave dedicated public-soak dual-interface client
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=flowweave
Group=flowweave
UMask=0077
ExecStart=/usr/local/libexec/flowweave-public-soak/flowweave-proxy client /etc/flowweave-public-soak/client.conf
Restart=on-failure
RestartSec=2s
TimeoutStopSec=15s
LimitNOFILE=65536
TasksMax=512
MemoryMax=1G

CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
